# debian-rootfs/files-tarragon/lib/systemd/system/srv.mount
[Unit]
Description=Data Partition File System

[Mount]
What=/dev/mmcblk1p5
Where=/srv
